# scripts/debian/usr/local/bin/dipicar
#!/usr/bin/python3
### BEGIN INIT INFO
# Provides:          dipicar
# Required-Start:    $all
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Starts the dipicar server
### END INIT INFO

from os import geteuid, chdir
from sys import argv
from subprocess import call
from re import match

INSTALLDIR = "/usr/lib/dipicar/"

class Dipicar:

  commands = []

  @staticmethod
  def config():
    call(["{}scripts/network/initUap0.py".format(INSTALLDIR)])
    call(["{}scripts/network/iptablesConfig.py".format(INSTALLDIR)])
  
  @staticmethod
  def version():
    call(["apt", "list", "dipicar"])
  
  @staticmethod
  def start():
    chdir(INSTALLDIR)
    call(["./dipicar_srv"])

  @staticmethod
  def fullStart():
    Dipicar.config()
    Dipicar.start()
  
  @staticmethod
  def genSSL():
    call(["{}scripts/installation/ssl_keys_gen.py".format(INSTALLDIR)])

Dipicar.commands = dir(Dipicar)[26:]

if(geteuid() != 0):
  print("You must run this script as root.")
  exit(0)

if (len(argv) > 2):
  print("This command only accept 1 argument.")
  exit(0)

if (len(argv) < 2):
  command = "version"
else:
  command = argv[1]

if command not in Dipicar.commands:
  print("You must use one of these commands: %s." %(", ".join(Dipicar.commands)))
  exit(0)

try:
  getattr(Dipicar, command)()
except KeyboardInterrupt:
  exit(0)
except Exception:
  exit(-1)